# task_cleanup.cron
0 * * * * /usr/local/bin/task_cleanup.sh